# config/supervisor/ubikom-server.conf
[program:ubikom-server]
command=/home/ubuntu/ubikom/ubikom-server --base-dir=/data/ubikom
