# config/supervisor/ubikom-web.conf
[program:ubikom-web]
command=/home/ubuntu/ubikom/ubikom-web --lookup-service-url=localhost:8825 --identity-service-url=localhost:8825 --cert-file=/etc/letsencrypt/live/alpha.ubikom.cc/fullchain.pem --key-file=/etc/letsencrypt/live/alpha.ubikom.cc/privkey.pem --ubikom-key-file=/home/ubuntu/ubikom/ubikom-web.key --ubikom-name=ubikom-web --notification-name=lgx --rate-limit-per-hour=20