# config/supervisor/ubikom-dump.conf
[program:ubikom-dump]
command=/home/ubuntu/ubikom/ubikom-dump --data-dir=/data/ubikom/dump
