# config/supervisor/ubikom-proxy.conf
[program:ubikom-proxy]
command=/home/ubuntu/ubikom/ubikom-proxy
