# config/supervisor-mail/ubikom-gateway.conf
[program:ubikom-gateway]
command=/home/ubuntu/ubikom/ubikom-gateway --key=/home/ubuntu/ubikom/gateway.key
